# meta-mellanox/meta-bluewhale-ast2500/recipes-phosphor/fans/files/obmc-mellanox-fand.service
[Unit]
Description=Mellanox Fan Controller

[Service]
Restart=always
RestartSec=60
ExecStart={sbindir}/obmc-mellanox-fand.sh

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
