# meta-mellanox/meta-bluewhale-ast2500/recipes-phosphor/mac-sync/files/obmc-mellanox-mac-syncd.service
[Unit]
Description=Mellanox MAC Sync Daemon

[Service]
ExecStart={sbindir}/obmc-mellanox-mac-syncd.sh

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
