# meta-mellanox/meta-bluewhale-ast2500/recipes-connectivity/ipmb-driver/files/obmc-mellanox-ipmbd.service
[Unit]
Description=Mellanox IPMB driver loader

[Service]
Restart=always
RestartSec=60
ExecStart={sbindir}/obmc-mellanox-ipmbd.sh

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
